# components/wallpaper/generate-wallpaper.timer
[Unit]
Description=Set a random wallpaper every 15 min

[Timer]
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
