# components/wallpaper/wallpaper-rotation.timer
[Unit]
Description=Set a random wallpaper every 15 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
