# components/battery-management/battery-level.timer
[Unit]
Description=Check battery level every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
